# llm-api-proxy.service
[Unit]
Description=LLM API Failover Proxy
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root
EnvironmentFile=/root/.llm-proxy.env
ExecStart=/usr/bin/python3 /root/llm-api-proxy.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
