# llm-api-proxy.service
[Unit]
Description=LLM API Failover Proxy
After=network.target

[Service]
Type=simple
# Change to your user (avoid running as root)
User=llmproxy
# Change to the directory containing llm-api-proxy.py
WorkingDirectory=/home/llmproxy
# Path to your .env file with API keys
EnvironmentFile=/home/llmproxy/.llm-proxy.env
# Adjust paths as needed: -c specifies config file location
ExecStart=/usr/bin/python3 /home/llmproxy/llm-api-proxy.py -c /home/llmproxy/.llm-proxy-config.json
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
